# initscripts/elivehealth
#!/bin/sh
### BEGIN INIT INFO
# Provides:          elivehealth
# Required-Start:    $local_fs $remote_fs $syslog
# Required-Stop:
# Should-Start:
# Should-Stop:
# X-Start-Before:
# X-Stop-After:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     false
# Short-Description: Configure elivehealth
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.
### END INIT INFO
# Policy Info:  http://wiki.debian.org/LSBInitScripts
# also check /etc/init.d/skeleton

# Author: Thanatermesis


# Defaults
#PIDDIR=/var/run/elivehealth
#elivehealthPID=$PIDDIR/elivehealth.pid
DESC="Configure elivehealth"
SCRIPTNAME=/etc/init.d/elivehealth


# Reads config file (will override defaults above)
[ -r /etc/default/elivehealth ] && . /etc/default/elivehealth

. /lib/lsb/init-functions


do_start(){
    local _CMDLINE
    _CMDLINE="$(cat /proc/cmdline)"
    if ! echo "${_CMDLINE}" | grep -qsE "(boot=live)"
    then
        exit 0
    fi
    if echo "${_CMDLINE}" | grep -q debug ; then
        set -x
    fi

    # oly if we are bash (not sh)
    #source /usr/lib/elive-tools/functions

    # make sure that the working console reads and outputs things correctly (otherwise tab not works on dialog's unless we do a "login" or "getty")
    #exec 0</dev/console 1>/dev/console 2>&1

    # if we had already the file is because we didn't shutdown correctly the system
    if test -e /var/tmp/elive-system-running ; then
        touch /var/tmp/elive-system-shutdown-forced
    fi

    touch /var/tmp/elive-system-running


    #log_daemon_msg "Starting $(basename elivehealth)"

    # Make sure we have our PIDDIR, even if it's on a tmpfs
    #install -o root -g root -m 755 -d $PIDDIR

    #log_progress_msg "daemonizing"

    #if ! start-stop-daemon --start --quiet --background --make-pidfile --pidfile $elivehealthPID --oknodo --exec /usr/sbin/elivehealthd -- -D
    #then
        #log_end_msg 1
        #exit 1
    #fi

    # turn off debugging if we have finished
    if echo "${_CMDLINE}" | grep -q debug ; then
        set +x
    fi

    #log_progress_msg "step 1 completed"
    true
    log_end_msg 0
}

do_stop(){
    #log_daemon_msg "Stopping elivehealth"
    #log_progress_msg "almost"
    rm -f /var/tmp/elive-system-running /var/tmp/elive-system-shutdown-forced

    #start-stop-daemon --stop --quiet --pidfile $elivehealthPID
    # Wait a little and remove stale PID file
    #if [ -f $elivehealthPID ] && ! ps h $(cat $elivehealthPID) > /dev/null
    #then
        ## Stale PID file (elivehealth was succesfully stopped),
        ## remove it (should be removed by elivehealth itself IMHO.)
    #rm -f $elivehealthPID
    #fi


    log_end_msg 0
}

case "$1" in
    "")
        echo "Error: elivehealth should be called with the 'start' argument." >&2
        ;;
    start)
        do_start
        ;;
    stop)
        do_stop
        ;;
    restart|reload|force-reload)
        "$0" stop
        "$0" start
        ;;
    *)
        echo "Usage: elivehealth [start|stop]" >&2
        exit 3
        ;;
esac

# vim: set filetype=sh :
